# scripts/backups/canary/pod-charlesreid1-canary.timer.j2
[Unit]
Description=Timer to run the pod-charlesreid1 backups canary

[Timer]
#OnCalendar=Sun *-*-* 9:03:00
OnCalendar=*-*-* 7:03:00

[Install]
WantedBy=timers.target
