# scripts/backups/pod-charlesreid1-backups-wikifiles.timer.j2
[Unit]
Description=Timer to back up the pod-charlesreid1 wiki database

[Timer]
OnCalendar=*-*-* 0/2:13:00

[Install]
WantedBy=timers.target

