# scripts/backups/pod-charlesreid1-backups-gitea.timer.j2
[Unit]
Description=Timer to back up pod-charlesreid1 gitea files

[Timer]
OnCalendar=Sun *-*-* 2:12:00
#OnCalendar=*-*-* 2:12:00

[Install]
WantedBy=timers.target
