# scripts/backups/pod-charlesreid1-backups-cleanolderthan.timer.j2
[Unit]
Description=Timer to clean files older than N days from the pod-charlesreid1 backups dir

[Timer]
OnCalendar=Sun *-*-* 2:28:00
#OnCalendar=*-*-* 2:28:00

[Install]
WantedBy=timers.target
