# scripts/backups/pod-charlesreid1-backups-aws.timer.j2
[Unit]
Description=Timer to copy the lastest pod-charlesreid1 backup to an S3 bucket

[Timer]
OnCalendar=Sat *-*-* 9:56:00

[Install]
WantedBy=timers.target
