# scripts/backups/pod-charlesreid1-backups-wikidb.timer.j2
[Unit]
Description=Timer to back up the pod-charlesreid1 wiki database

[Timer]
OnCalendar=Sat *-*-* 9:03:00

[Install]
WantedBy=timers.target
